# myapp.service
[Unit]
Description=My Awesome App
Documentation=https://mycompany.com/myawesomeapp

[Service]
ExecStart=/usr/bin/myapp --level=awesome
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target

#sudo mv myapp.service /etc/systemd/system/
#sudo systemctl daemon-reload
#sudo systemctl enable myapp
#sudo systemctl start myapp
#sudo systemctl status myapp --no-pager
